# roles/common/files/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.

setterm --blank 1

exit 0
